# files/shared/etc/systemd/system/custom.service
[Unit]
Description=Custom Service
After=network.target

[Service]
ExecStart=/usr/libexec/CUSTOM
Restart=always

[Install]
WantedBy=multi-user.target
