# config/files/usr/etc/profile.d/myos-neofetch.sh
alias neofetch='/usr/bin/neofetch --source /usr/share/ublue-os/myos/logo.txt --ascii_colors 12 7 13 --config /usr/share/ublue-os/myos/neofetch.conf --package_separate off'
